# haproxy-agent/deb/haproxy-agent.service
[Unit]
Description=Haproxy Agent

[Service]
ExecStart=/usr/bin/haproxy-agent

[Install]
WantedBy=multi-user.target